# root/etc/s6-overlay/s6-rc.d/init-kometa-config/up
/etc/s6-overlay/s6-rc.d/init-kometa-config/run
